# deploy/systemd/finpilot-fx.timer
[Unit]
Description=FINPILOT — обновление курсов ЦБ по будням (12:30)

[Timer]
OnCalendar=Mon-Fri 12:30:00
Persistent=true

[Install]
WantedBy=timers.target
